# scripts/ocp/monitor-ocp.service
[Unit]
Description=OCP Dashboard Monitoring Script
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/bin/python3 /opt/monitoring/monitor_ocp.py
User=root
Environment="DASHBOARD_API_URL=https://overview/api/components"
